# deploy/dockpilot-agent.service
[Unit]
Description=DockPilot Agent
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=root
Environment=DOCKPILOT_SERVER_URL=http://127.0.0.1:8080
Environment=DOCKPILOT_REGISTRATION_TOKEN=change-me-registration-token
Environment=DOCKPILOT_STATE_PATH=/var/lib/dockpilot-agent/agent.json
Environment=DOCKPILOT_COMPOSE_DIRS=/opt,/srv,/var/www
Environment=TZ=Asia/Shanghai
EnvironmentFile=-/etc/dockpilot/agent.env
ExecStart=/opt/dockpilot-agent/dockpilot-agent
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
